# config/.config/autostart/compton.desktop
[Desktop Entry]
Name=compton
Comment=X Composition Manager
Exec=compton -c
Terminal=false
Type=Application
NoDisplay=true
NotShowIn=Unity
