# config/.config/autostart/numlockx.desktop
[Desktop Entry]
Name=numlockx
Comment=Numlock at startup
Exec=numlockx
Terminal=false
Type=Application
NoDisplay=true
NotShowIn=Unity;LXDE;
